# services/user-service/kafka_consumers.conf
[program:EXAMPLE]
command=php /app/backend/artisan kafka:consume:EXAMPLE
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/kafka_EXAMPLE.log